# platform/orange-pi/rootfs/etc/systemd/system/jackd.service
[Unit]
Description=JACK Audio Connection Kit
After=systemd-udev-settle.service sound.target
Wants=sound.target
Before=openrig.service

[Service]
Type=simple
User=root
Group=audio

# Stability config for Rockchip RK3588 xHCI + USB audio: 256 frames/period, 3 periods = ~16ms @ 48kHz.
# Scarlett Gen 4 is locked to 48kHz by the mainline scarlett-gen2 driver — lower rates cause
# "Error initialising Scarlett Gen 4 Mixer Driver: -71" (EPROTO). Small buffers on RK3588 xHCI
# have been correlated with host controller resets and interface disconnects (red LED).
# Waits up to 10s for any USB audio interface; exits cleanly if none found.
Environment=JACK_NO_AUDIO_RESERVATION=1

# systemd ignores PAM limits.conf, so RT scheduling must be granted here.
LimitRTPRIO=95
LimitMEMLOCK=infinity

# Pin xHCI USB host IRQs to a big A76 core (CPU 4 = mask 0x10) to reduce audio jitter.
# RK3588S topology: CPUs 0-3 = A55 (little), CPUs 4-7 = A76 (big).
ExecStartPre=/bin/sh -c '\
  for irq in $(grep -E "xhci-hcd" /proc/interrupts 2>/dev/null | awk -F: "{print \$1}" | tr -d " "); do \
    echo 10 > /proc/irq/$irq/smp_affinity 2>/dev/null || true; \
  done; true'

ExecStart=/bin/sh -c '\
  for i in $(seq 1 10); do \
    CARD=$(aplay -l 2>/dev/null | grep "USB Audio" | head -1 | sed "s/card \([0-9]*\):.*/\1/"); \
    if [ -n "$CARD" ]; then \
      echo "Found USB audio card $CARD, starting JACK..."; \
      exec /usr/bin/jackd -d alsa -d hw:$CARD -r 48000 -p 256 -n 3; \
    fi; \
    sleep 1; \
  done; \
  echo "No USB audio interface found after 10s, JACK disabled (OpenRig will use ALSA)"'

Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
