# platform/orange-pi/rootfs/etc/systemd/system/openrig.service
[Unit]
Description=OpenRig Guitar Pedalboard
Documentation=https://github.com/jpfaria/OpenRig
After=jackd.service weston.service
Wants=jackd.service
Requires=weston.service

[Service]
Type=simple
User=root
Group=root

Environment=WAYLAND_DISPLAY=wayland-1
Environment=XDG_RUNTIME_DIR=/run/weston
Environment=RUST_LOG=warn
Environment=JACK_NO_AUDIO_RESERVATION=1
# Native libs (NAM + LV2 bundles) are shipped by the openrig .deb under
# /usr/lib/openrig/libs/ (per FHS). Some LV2 plugins have bundled side
# libraries that must also be resolvable by the dynamic linker, so both
# directories are added here. /usr/share/openrig/libs -> /usr/lib/openrig/libs
# is also symlinked in customize-image.sh for the OpenRig data-root lookup.
Environment=LD_LIBRARY_PATH=/usr/lib/openrig/libs/nam/linux-aarch64:/usr/lib/openrig/libs/lv2/linux-aarch64
Environment=GDK_BACKEND=wayland

# Wait for weston to create the Wayland socket before starting.
ExecStartPre=/bin/sh -c 'for i in $(seq 1 15); do [ -S /run/weston/wayland-1 ] && exit 0; sleep 1; done; exit 1'

# Wait for jackd to be fully up before starting OpenRig.
#
# CRITICAL: jackd.service is Type=simple wrapped in a bash loop that waits
# up to 10s for a USB audio interface to appear before exec'ing jackd. That
# means systemd marks jackd.service as "started" the moment the bash
# wrapper forks — NOT when jackd itself has actually opened ALSA and is
# serving clients. If OpenRig starts in that window, jack_server_is_running()
# in infra-cpal returns false (no /dev/shm/jack_default_*_0 yet), OpenRig
# falls into the CPAL/ALSA path, and supported_input_configs() probes the
# Scarlett via ALSA at the exact moment jackd is also about to open it.
# On RK3588 xHCI this double-open causes the Scarlett Gen 4 to fall off the
# USB bus (red LED, card disappears until replug or reboot).
#
# Wait up to 20s for the JACK server socket to appear. If it never does
# (no USB audio present), give up and start anyway — OpenRig will fall
# back to silent/ALSA mode and at least the UI will come up.
ExecStartPre=/bin/sh -c 'for i in $(seq 1 20); do ls /dev/shm/jack_default_*_0 2>/dev/null | grep -q . && exit 0; sleep 1; done; echo "WARN: JACK socket not found after 20s, starting OpenRig anyway"; exit 0'

ExecStart=/usr/bin/openrig /home/openrig/project.yaml --auto-save
Restart=always
RestartSec=5
TimeoutStartSec=30

[Install]
WantedBy=multi-user.target
