# platform/orange-pi/rootfs/etc/systemd/system/openrig-irq-affinity.service
[Unit]
Description=Pin xHCI USB IRQs to big A76 core for audio jitter reduction
After=systemd-udev-settle.service
Before=openrig.service

[Service]
Type=oneshot
RemainAfterExit=yes

# Pin xHCI USB host IRQs to a big A76 core (CPU 4 = mask 0x10) to reduce audio jitter.
# RK3588S topology: CPUs 0-3 = A55 (little), CPUs 4-7 = A76 (big).
ExecStart=/bin/sh -c '\
  for irq in $(grep -E "xhci-hcd" /proc/interrupts 2>/dev/null | awk -F: "{print \$1}" | tr -d " "); do \
    echo 10 > /proc/irq/$irq/smp_affinity 2>/dev/null || true; \
  done; true'

[Install]
WantedBy=multi-user.target
