# platform/orange-pi/rootfs/etc/systemd/system/openrig.service.d/cpuaffinity.conf
[Service]
# Pin the openrig process (and every thread it creates: Slint UI, llvmpipe
# software GPU renderer for Mesa, spectrum-analyzer poll_stream workers,
# smithay-clipboard, etc) to the little A55 cores 0-3.
#
# The DSP worker thread widens its own affinity mask to the big A76 cores
# (4-7) via sched_setaffinity in infra-cpal. sched_setaffinity may widen
# the mask — the CPUAffinity= directive above only sets the INITIAL mask,
# it does not impose a cgroup cpuset restriction.
#
# jackd is pinned to the big cores separately via pre_exec + sched_setaffinity
# in LiveJackBackend::spawn (infra-cpal/src/jack_supervisor/live_backend.rs).
#
# Rationale: the kernel cmdline isolcpus=4-7 is SOFT on modern kernels —
# it only excludes the big cores from default scheduler load balancing; it
# does NOT prevent threads from being spawned there. Without this pin, the
# llvmpipe GPU renderer spawns worker threads that can land on big cores 5,
# 6 and 7, competing for CPU with the JACK RT callback and the DSP worker.
# The symptom was audible clicks when moving the mouse over UI blocks (the
# spectrum-analyzer poll_stream + Slint repaint path).
CPUAffinity=0-3
