# orange-pi/rootfs/etc/systemd/system/jackd.service
[Unit]
Description=JACK Audio Connection Kit
After=systemd-udev-settle.service sound.target
Wants=sound.target
Before=openrig.service

[Service]
Type=simple
User=openrig
Group=openrig
SupplementaryGroups=audio

# Wait for Teyun Q-26 to be enumerated (max 10s)
ExecStartPre=/bin/sh -c 'for i in $(seq 1 10); do aplay -l 2>/dev/null | grep -q Q26 && break; sleep 1; done'

# Low-latency config: 128 frames/period, 2 periods = ~5ms @ 48kHz
Environment=JACK_NO_AUDIO_RESERVATION=1
ExecStart=/usr/bin/jackd -d alsa -d hw:Q26 -r 48000 -p 128 -n 2

Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
