# orange-pi/rootfs/etc/systemd/system/jackd.service
[Unit]
Description=JACK Audio Connection Kit
After=systemd-udev-settle.service sound.target
Wants=sound.target
Before=openrig.service

[Service]
Type=simple
User=openrig
Group=openrig
SupplementaryGroups=audio

# Low-latency config: 128 frames/period, 2 periods = ~5ms @ 48kHz
# Waits up to 10s for any USB audio interface; exits cleanly if none found.
Environment=JACK_NO_AUDIO_RESERVATION=1
ExecStart=/bin/sh -c '\
  for i in $(seq 1 10); do \
    CARD=$(aplay -l 2>/dev/null | awk "/USB Audio/ {match(\$0, /card ([0-9]+)/, m); if (m[1]!=\"\") {print m[1]; exit}}"); \
    if [ -n "$CARD" ]; then \
      echo "Found USB audio card $CARD, starting JACK..."; \
      exec /usr/bin/jackd -d alsa -d hw:$CARD -r 48000 -p 128 -n 2; \
    fi; \
    sleep 1; \
  done; \
  echo "No USB audio interface found after 10s, JACK disabled (OpenRig will use ALSA)"'

Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
