# orange-pi/rootfs/etc/systemd/system/openrig.service
[Unit]
Description=OpenRig Guitar Pedalboard
Documentation=https://github.com/jpfaria/OpenRig
After=jackd.service weston.service
Wants=jackd.service
Requires=weston.service

[Service]
Type=simple
User=openrig
Group=openrig
SupplementaryGroups=audio video

Environment=WAYLAND_DISPLAY=wayland-1
Environment=XDG_RUNTIME_DIR=/run/user/0
Environment=RUST_LOG=warn
Environment=JACK_NO_AUDIO_RESERVATION=1
Environment=LD_LIBRARY_PATH=/usr/local/share/openrig/libs/nam/linux-aarch64

ExecStart=/usr/local/bin/openrig
Restart=on-failure
RestartSec=5
TimeoutStartSec=30

[Install]
WantedBy=multi-user.target
